# documentacion/sistema_vigilancia.service
[Unit]
Description=Sistema de Vigilancia Aut贸nomo SADA
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/protocolo_deteccion
ExecStart=/usr/bin/python3 /home/pi/protocolo_deteccion/sistema_vigilancia_autonomo.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Configuraci贸n para Raspberry Pi
Environment=PYTHONPATH=/home/pi/protocolo_deteccion
Environment=DISPLAY=:0

# Configuraci贸n de recursos
MemoryMax=512M
CPUQuota=80%

# Configuraci贸n de seguridad
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/protocolo_deteccion

[Install]
WantedBy=multi-user.target
